# rootfs/etc/s6/sdrplay_api/run
#!/bin/bash

exec /usr/bin/sdrplay_apiService